# s6/s6-overlay/s6-rc.d/telly-log-prepare/up
if { mkdir -p /var/log/telly }
if { chown nobody:nogroup /var/log/telly }
chmod 02755 /var/log/telly