# s6/s6-overlay/s6-rc.d/telly/run
#!/command/execlineb -P

exec /usr/bin/telly 2>&1