# s6/s6-overlay/s6-rc.d/init-envfile/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

source /etc/env

echo "
RUNNING ENVIRONMENT
-----------------------------------------------
PUID                   | $PUID
PGID                   | $PGID
VERSION                | $VERSION
STREAMS                | $STREAMS
FFMPEG                 | $FFMPEG
BASE                   | $BASE
M3U                    | $M3U
EPG                    | $EPG
FILTER                 | $FILTER
PERSISTENCE            | $PERSISTENCE
DEVICE_AUTH            | $DEVICE_AUTH
DEVICE_ID              | $DEVICE_ID
DEVICE_UUID            | $DEVICE_UUID
DEVICE_FIRMWARE_NAME   | $DEVICE_FIRMWARE_NAME
DEVICE_FIRMWARE_VERSION| $DEVICE_FIRMWARE_VERSION
DEVICE_FRIENDLY_NAME   | $DEVICE_FRIENDLY_NAME
DEVICE_MANUFACTURER    | $DEVICE_MANUFACTURER
DEVICE_MODEL_NUMBER    | $DEVICE_MODEL_NUMBER
SSDP                   | $SSDP
STARTING_CHANNEL       | $STARTING_CHANNEL
XMLTV_CHANNELS         | $XMLTV_CHANNELS
LOG_LEVEL              | $LOG_LEVEL
LOG_REQUESTS           | $LOG_REQUESTS
BASE_ADDRESS           | $BASE_ADDRESS
LISTEN_ADDRESS         | $LISTEN_ADDRESS
SCHEDULES_DIRECT_USERNAME | $SCHEDULES_DIRECT_USERNAME
SCHEDULES_DIRECT_PASSWORD | $SCHEDULES_DIRECT_PASSWORD
SOURCE1_NAME           | $SOURCE1_NAME
SOURCE1_PROVIDER       | $SOURCE1_PROVIDER
SOURCE1_USERNAME       | $SOURCE1_USERNAME
SOURCE1_PASSWORD       | $SOURCE1_PASSWORD
SOURCE1_FILTER         | $SOURCE1_FILTER
SOURCE1_FILTER_KEY     | $SOURCE1_FILTER_KEY
SOURCE1_FILTER_RAW     | $SOURCE1_FILTER_RAW
SOURCE1_SORT           | $SOURCE1_SORT
SOURCE2_NAME           | $SOURCE2_NAME
SOURCE2_PROVIDER       | $SOURCE2_PROVIDER
SOURCE2_USERNAME       | $SOURCE2_USERNAME
SOURCE2_PASSWORD       | $SOURCE2_PASSWORD
SOURCE3_PROVIDER       | $SOURCE3_PROVIDER
SOURCE3_M3U            | $SOURCE3_M3U
SOURCE3_EPG            | $SOURCE3_EPG
-----------------------------------------------
"